# polymarket-bot.service
[Unit]
Description=Polymarket Final Bot
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=3
 
[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/polytest-final
ExecStart=/home/ubuntu/polytest-final/venv/bin/python main.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
EnvironmentFile=/home/ubuntu/polytest-final/.env
 
[Install]
WantedBy=multi-user.target